# agent/nccl-agent.service.template
[Unit]
Description=NCCL Collective Notify Forwarding Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
ExecStart=__AGENT_BIN__
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
